# 10-lisenet-exam/roles/sample-apache/templates/restart-httpd.path.j2
[Install]
WantedBy=httpd.service
   
[Path]
PathModified={{ systemd_path_to_monitor }}
