# 10-lisenet-exam/roles/sample-apache/files/restart-httpd.service
[Service]
ExecStart=systemctl restart httpd.service
Type=oneshot
